# source/etc/crontab
#!/usr/bin/crontab 

* * * * * $HOME/bin/readfreq
2,17,32,47 * * * * $HOME/bin/readmidc
51 * * * * $HOME/bin/nextmidc
4,9,14,19,24,29,34,39,44,49,54,59 * * * * $HOME/bin/feederbid
0,5,10,15,20,25,30,35,40,45,50,55 * * * * $HOME/bin/clearmkt
3,13,23,33,43,53 * * * * $HOME/bin/getweather
53,54,55,56 * * * * $HOME/bin/getmetar
58 * * * * $HOME/bin/convertmetar
0 * * * * $HOME/bin/gfacheck
1 1 * * * $HOME/bin/getgoodWdata
